# src/charon/plugins/nm/gnome/debian/network-manager-strongswan.postinst
#!/bin/sh

set -e

case "$1" in
	configure)
		if [ -x "/etc/init.d/dbus" ]; then
			if [ -x /usr/sbin/invoke-rc.d ]; then
				invoke-rc.d dbus force-reload || true
			else
				/etc/init.d/dbus force-reload || true
			fi
		fi
        ;;
esac

#DEBHELPER#

exit 0
