# debian/strongswan.postinst
#!/bin/sh
set -e

cp -f /etc/ssl/openssl.cnf /etc/ssl/openssl.cnf.orig
systemctl unmask strongswan-starter.service
systemctl enable strongswan-starter.service

#DEBHELPER#

exit 0
