# testing/tests/ikev2/dhcp-static-client-id/hosts/venus/etc/init.d/dhcpd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="start stop"

depend() {
	need net
	need logger
}

start() {
	ebegin "Starting DHCP server"
	start-stop-daemon --start --quiet --exec /usr/sbin/dhcpd
	eend $?
}

stop() {
	ebegin "Stopping DHCP server"
	start-stop-daemon --stop --quiet --pidfile /var/run/dhcpd.pid
	rm -f /var/state/dhcp/dhcpd.leases
	touch /var/state/dhcp/dhcpd.leases
	eend $?
}
